# infrastructure/systemd/alpha-engine-daemon.service
[Unit]
Description=Alpha Engine Intraday Trading Daemon
After=ibgateway.service
Requires=ibgateway.service

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/alpha-engine
ExecStartPre=/home/ec2-user/alpha-engine/infrastructure/wait-for-ibgateway.sh
ExecStart=/home/ec2-user/alpha-engine/.venv/bin/python -m executor.daemon
ExecStartPost=/bin/bash /home/ec2-user/alpha-engine/infrastructure/emit-heartbeat.sh executor-daemon-start
ExecStopPost=/bin/sleep 5
TimeoutStopSec=30
Restart=on-failure
RestartSec=30
StandardOutput=append:/var/log/daemon.log
StandardError=append:/var/log/daemon.log
Environment=PYTHONPATH=/home/ec2-user/alpha-engine
Environment=AWS_REGION=us-east-1
Environment=ALPHA_ENGINE_JSON_LOGS=1
Environment=FLOW_DOCTOR_ENABLED=1
Environment=ALPHA_ENGINE_DEPLOYED=1
# Secrets — NO EnvironmentFile (.env deprecated, #890). Every secret the
# daemon and its transitive imports touch is now obtained at runtime via
# nousergon_lib.secrets.get_secret() (SSM-backed, lib v0.62.0):
#   - EMAIL_SENDER / EMAIL_RECIPIENTS / GMAIL_APP_PASSWORD → email_sender.py
#   - TELEGRAM_BOT_TOKEN / TELEGRAM_CHAT_ID → telegram.py
#   - flow-doctor's ${EMAIL_*}/${GMAIL_APP_PASSWORD}/${TELEGRAM_*}/
#     ${FLOW_DOCTOR_GITHUB_TOKEN} → seeded into env from get_secret() by
#     logging.py::_seed_flow_doctor_secrets() at setup_logging() import time,
#     which is authoritative under ALPHA_ENGINE_DEPLOYED=1 (set above).
#   - ANTHROPIC_API_KEY → no longer referenced (executor LLM path removed
#     2026-05-25; "no LLM in executor, period").
# This mirrors alpha-engine-morning.service, which already runs the same
# executor codebase via systemd with no EnvironmentFile. The
# tests/test_no_secret_environ_reads.py regression guard fails CI if any of
# these are re-introduced as os.environ reads. The historical 2026-05-15
# crash-loop (flow-doctor ${EMAIL_SENDER} unresolved after the ssm_secrets
# import shim was deleted in #182) is resolved by the get_secret() seeding.

[Install]
WantedBy=multi-user.target
