# infrastructure/systemd/alpha-engine-daemon.timer
[Unit]
Description=Start Alpha Engine Intraday Daemon before market open

[Timer]
# Start at 9:29 AM ET (6:29 AM PT), 1 min before market open.
# Safety net — daemon also starts on boot via multi-user.target.
OnCalendar=Mon..Fri 09:29 America/New_York
Persistent=false

[Install]
WantedBy=timers.target
