# infrastructure/systemd/alpha-engine-daemon.service
[Unit]
Description=Alpha Engine Intraday Trading Daemon
After=ibgateway.service
Requires=ibgateway.service

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/alpha-engine
ExecStartPre=/home/ec2-user/alpha-engine/infrastructure/wait-for-ibgateway.sh
ExecStart=/home/ec2-user/alpha-engine/.venv/bin/python -m executor.daemon
ExecStartPost=/bin/bash /home/ec2-user/alpha-engine/infrastructure/emit-heartbeat.sh executor-daemon-start
ExecStopPost=/bin/sleep 5
TimeoutStopSec=30
Restart=on-failure
RestartSec=30
StandardOutput=append:/var/log/daemon.log
StandardError=append:/var/log/daemon.log
Environment=PYTHONPATH=/home/ec2-user/alpha-engine
Environment=AWS_REGION=us-east-1
Environment=ALPHA_ENGINE_JSON_LOGS=1
Environment=FLOW_DOCTOR_ENABLED=1
# Secrets (EMAIL_SENDER, GMAIL_APP_PASSWORD, ANTHROPIC_API_KEY, ...) — the
# RunDaemon SSM step launches this unit via systemd, which (unlike the SF's
# RunMorningPlanner command) does not `source` the env file. PR 9g (#182)
# deleted the ssm_secrets import-time shim that previously populated these
# for the daemon, so flow-doctor's ${EMAIL_SENDER} resolution crash-looped
# the daemon (2026-05-15). Mirror the morning-planner path until the
# .env-deprecation arc's SSM substrate fully covers systemd entrypoints.
# Leading `-` tolerates the file's eventual removal (no hard failure).
EnvironmentFile=-/home/ec2-user/.alpha-engine.env

[Install]
WantedBy=multi-user.target
