# infrastructure/systemd/alpha-engine-daemon.service
[Unit]
Description=Alpha Engine Intraday Trading Daemon
After=ibgateway.service
Requires=ibgateway.service

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/alpha-engine
EnvironmentFile=/home/ec2-user/.alpha-engine.env
ExecStart=/home/ec2-user/alpha-engine/.venv/bin/python -m executor.daemon
Restart=on-failure
RestartSec=30
StandardOutput=append:/var/log/daemon.log
StandardError=append:/var/log/daemon.log
Environment=PYTHONPATH=/home/ec2-user/alpha-engine

[Install]
WantedBy=multi-user.target
