# infrastructure/systemd/alpha-engine-eod.timer
[Unit]
Description=Trigger Alpha Engine EOD Reconciliation at 1:20 PM PT

[Timer]
# Run after 4:15 PM ET (1:15 PM PT) to ensure all 15-min delayed
# market data has arrived before computing NAV and daily P&L.
OnCalendar=Mon..Fri 13:20 America/Los_Angeles
Persistent=false

[Install]
WantedBy=timers.target
